# etc_systemd_system/clear-matrix.service
#[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/clear-matrix.sh

[Install]
WantedBy=default.target
