# etc_systemd_system/pyaudio-spectre.service
#[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/pyaudio-spectre.sh

[Install]
WantedBy=default.target
